# bin/impress.service
[Unit]
Description=Impress Application Server
After=network.target
[Service]
ExecStart=/bin/sh -c 'impress start'
ExecStop=/bin/sh -c 'impress stop'
ExecReload=/bin/sh -c 'impress restart'
Type=forking
OOMScoreAdjust=-1000
CPUSchedulingPolicy=fifo
CPUSchedulingPriority=99
IOSchedulingClass=realtime
IOSchedulingPriority=0
[Install]
WantedBy=multi-user.target
