# bin/impress
#!/bin/sh

# chkconfig: 2345 90 10
# description: Impress Application Server
# processname: impress

### BEGIN INIT INFO
# Provides:          impress
# Required-Start:    $syslog
# Required-Stop:     $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Impress Application Server
# Description:       Impress Application Server provides universal platform for Node.js applications
### END INIT INFO

case $1 in

  start)
    impress start
    ;;
  stop)
    impress stop
    ;;
  restart|reload)
    impress restart
    ;;
  update)
    impress update
    ;;
  status)
    impress status
    ;;
  *)
    impress help
    ;;

esac
